# service/car-tracker.timer
[Unit]
Description="Pull car tracker pages every day"

[Timer]
OnCalendar=daily
Persitstent=true

[Install]
WantedBy=timers.target